# asl3-update-nodelist.service
[Unit]
Description=Update ASL3 on-disk nodelist db

[Service]
Type=oneshot
User=asterisk
Group=asterisk
ExecStart=/usr/bin/asl3-update-nodelist
